# systemd/docwain-vllm-smart.service
[Unit]
Description=DocWain vLLM Smart Instance (27B)
After=network.target
Wants=network-online.target
# Smart and Fast conflict on GPU — only run one at a time
Conflicts=docwain-vllm-fast.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/PycharmProjects/DocWain
Environment="PATH=/home/ubuntu/PycharmProjects/DocWain/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="CUDA_VISIBLE_DEVICES=0"
ExecStart=/home/ubuntu/PycharmProjects/DocWain/.venv/bin/python -m vllm.entrypoints.openai.api_server \
    --model Qwen/Qwen3.5-27B \
    --served-model-name docwain-smart \
    --port 8200 \
    --host 0.0.0.0 \
    --dtype bfloat16 \
    --max-model-len 32768 \
    --gpu-memory-utilization 0.85 \
    --enable-chunked-prefill \
    --tensor-parallel-size 1
Restart=on-failure
RestartSec=10
TimeoutStartSec=600
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
