# systemd/docwain-vllm-fast.service
[Unit]
Description=DocWain vLLM Fast Instance (14B)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/PycharmProjects/DocWain
Environment="PATH=/home/ubuntu/PycharmProjects/DocWain/.venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="CUDA_VISIBLE_DEVICES=0"
ExecStart=/home/ubuntu/PycharmProjects/DocWain/.venv/bin/python -m vllm.entrypoints.openai.api_server \
    --model models/DocWain-14B-v2-AWQ \
    --served-model-name docwain \
    --port 8100 \
    --host 0.0.0.0 \
    --dtype auto \
    --quantization awq \
    --max-model-len 32768 \
    --gpu-memory-utilization 0.90 \
    --enable-prefix-caching \
    --enable-chunked-prefill \
    --tensor-parallel-size 1
Restart=on-failure
RestartSec=10
TimeoutStartSec=300
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
