# systemd/docwain-standalone.service
[Unit]
Description=DocWain Standalone API
After=network.target mongod.service
Wants=mongod.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/PycharmProjects/DocWain
Environment=STANDALONE_PORT=8400
Environment=VLLM_BASE_URL=http://localhost:8100/v1
Environment=VLLM_MODEL_NAME=docwain
Environment=STANDALONE_MONGODB_URI=mongodb://localhost:27017
Environment=STANDALONE_MONGODB_DB=docwain_standalone
EnvironmentFile=-/home/ubuntu/PycharmProjects/DocWain/.env.standalone
ExecStart=/home/ubuntu/PycharmProjects/DocWain/.venv/bin/python -m standalone
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
